# deployment/ica-schedule.service
[Unit]
Description=Periodical FastAPI Application
Documentation=https://github.com/your-repo/ica-schedule
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/Periodical

# Environment variables
Environment="SECRET_KEY=CHANGE_ME_TO_RANDOM_SECRET"
Environment="PRODUCTION=true"
Environment="PYTHONUNBUFFERED=1"

# Load additional env vars from file (optional)
# EnvironmentFile=/opt/Periodical/.env

# Start command - Använd venv för isolerad miljö
ExecStart=/opt/Periodical/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=0

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/Periodical/app/database
ReadWritePaths=/opt/Periodical/logs

# Resource limits
LimitNOFILE=4096
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ica-schedule

[Install]
WantedBy=multi-user.target
